# ide50-2015/files/etc/cron.d/appliance50_rsnapshot
*/10 * * * * root /usr/bin/nice /usr/bin/rsnapshot minutely
2 * * * * root /usr/bin/nice /usr/bin/rsnapshot hourly
4 0 * * * root /usr/bin/nice /usr/bin/rsnapshot daily
6 0 * * 0 root /usr/bin/nice /usr/bin/rsnapshot weekly
